# scripts/avian-record.timer
[Unit]
Description=Run Avian Recorder every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target