# scripts/avian-update.service
[Unit]
Description=Update Avian Recorder

[Service]
Type=oneshot
ExecStart=/home/birdpi/avian-acoustic-monitoring/scripts/update_and_rebuild.sh